# supervisord.conf
[inet_http_server]
port=0.0.0.0:9001
username=faleiro
password=faleiro

[unix_http_server]
file=/tmp/faleiro-supervisor.sock

[supervisord]
logfile=/var/log/faleiro/faleiro-supervisord.log
logfile_maxbytes=10MB
logfile_backups=10
loglevel=warn
pidfile=/var/run/faleiro-supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=root
childlogdir=/var/log/faleiro/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
prompt = faleiro
serverurl=unix:///tmp/faleiro-supervisor.sock


[program:faleiro-0]
command=./run.sh
stderr_logfile = /var/log/faleiro/stderr.log
stdout_logfile = /var/log/faleiro/stdout.log
